# homelight-api.service
[Unit]
Description=Homelight-API

[Service]
PIDFile=/tmp/homelight-api.pid
User=Pi
Restart=always
RestartSec=1
KillSignal=SIGQUIT
ExecStart=/home/pi/Homelight-API/server.js

[Install]
WantedBy=multi-user.target
